# pivotalzibot.conf
[program:trackerzibot]
command=sh /usr/local/share/trackerzibot/start_trackerzibot.sh
startsecs=20
stdout_logfile=/var/log/trackerzibot.log
redirect_stderr=true
