# pivotalzibot.conf
[program:pivotalzibot]
command=sh /usr/local/share/pivotalzibot/start_pivotalzibot.sh
startsecs=20
stdout_logfile=/var/log/pivotalzibot.log
redirect_stderr=true
